# ancillary scripts/ninja-store.conf
#!upstart
# Upstart script belongs in /etc/init/comet-store.conf
# This allows jenkins to call "sudo start comet-store" and "sudo stop comet-store"

description "comet-store"
author      "Stuart"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

# dont respawn while developing
# respawn                # restart when job dies
# respawn limit 5 60     # give up restart after 5 respawns in 60 seconds

script
  exec  node "/var/lib/jenkins/workspace/Comet Store/Router.js" >> /home/stuart/simple.log 2>&1
end script
